# nubis/puppet/files/rc.local
#!/bin/sh

LOGGER="logger --stderr --priority local7.info --tag nubis-startup"

if [ -d /etc/nubis.d ]; then
  run-parts /etc/nubis.d 2>&1 | $LOGGER
fi

exit 0
